# scripts/add-ipv6.service
[Unit]
Description=Check IPV6 service
After=network.target

[Service]
ExecStart=/usr/local/bin/add-ipv6.sh

[Install]
WantedBy=multi-user.target
